# deploy/wealthfolio-backend.service
[Unit]
Description=Wealthfolio Backend (FastAPI/Uvicorn)
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=%h/wealthfolio/backend
ExecStart=%h/wealthfolio/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

# Logging
StandardOutput=append:%h/wealthfolio/backend/server.log
StandardError=append:%h/wealthfolio/backend/server.log

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h/wealthfolio/backend

[Install]
WantedBy=multi-user.target
